# backtalk-br2-external/overlay/etc/systemd/system/bthelper@.service
[Unit]
Description=Raspberry Pi bluetooth helper
Requires=bluetooth.service
After=bluetooth.service

[Service]
Type=oneshot
ExecStart=/usr/libexec/bthelper.sh %I
Restart=on-failure
RemainAfterExit=yes
