# backtalk-br2-external/overlay/etc/systemd/system/asterisk.service.d/asterisk-bluetooth.conf
[Service]
ExecStartPre=/usr/libexec/asterisk-bluetooth.rb
